# sources/meta-fsl-bsp-release/imx/meta-bsp/recipes-core/initscripts/arm/imx/rtc-sync.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides: resize
# Required-Start:    mountkernfs
# Required-Stop:     mountkernfs
# Default-Start:     S
# Default-Stop:
### END INIT INFO

echo "finish" > /etc/flag

sleep 10

hwclock --hctosys

: exit 0

